# src/supervisor_conf.conf
[program:yotrawat_bot]
command=/usr/bin/python3 /yotrawat/src/bot_app.py
directory=/yotrawat/src/
user=Shadiv
autostart=true
autorestart=true
redirect_stderr=true
